# docker/supervisor-app.conf
[supervisord]
[program:app-uwsgi]
command = /usr/sbin/uwsgi --ini /home/hospital/app/docker/uwsgi.ini

[program:app-nginx]
command = /usr/sbin/nginx -g 'pid /tmp/nginx.pid; daemon off;'