# docker/supervisor-app.conf
[supervisord]
[program:app-uwsgi]
command = /bin/ash -c "$APP_DIR/docker/uwsgi.run"

[program:app-nginx]
command = /bin/ash -c "$APP_DIR/docker/nginx.run"